# receiver/audiocast.service
[Unit]
Description=audiocast

[Service]
Type=simple
ExecStart=/usr/bin/audiocast
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
